# install/systemd/omarchy-limine-snapshot.service
[Unit]
Description=Update Limine configuration when snapshots change
After=snapper-cleanup.service

[Service]
Type=oneshot
User=root
ExecStart=omarchy-limine-snapshot-hook
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/jon/.local/share/omarchy/bin
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=10s

# Security hardening (relaxed for symlink access)
CapabilityBoundingSet=CAP_SYS_ADMIN
LockPersonality=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHome=no
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ReadWritePaths=/tmp /boot /boot/EFI /home
RemoveIPC=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictSUIDSGID=yes
NoNewPrivileges=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service @mount

[Install]
WantedBy=multi-user.target